# scripts/bootup.service
[Unit]
Description=Setting Bootup Scripts and Others

[Service]
ExecStart = /usr/bin/bootup_script.sh

[Install]
WantedBy=multi-user.target

